# kubeinstall/build/node/services/terrakube-node.service
[Unit]
Description=Terrakube - Joins Kubernetes node to master
After=network.target

[Service]
PIDFile=/run/terrakube-node.pid
ExecStart=/usr/sbin/terrakube-node
ExecStop=-/sbin/start-stop-daemon --quiet --stop --retry QUIT/5 --pidfile /run/terrakube-node.pid
TimeoutStopSec=5
KillMode=mixed
Restart=always
StartLimitInterval=0
RestartSec=10


[Install]
WantedBy=multi-user.target
